# services/ffmpeg.ini
[program:ffmpeg]
command=sh -c 'set -x; if [ "$DEBUG" == "true" ]; then /usr/bin/ffmpeg -f x11grab -video_size 1920x1080 -i :99 -r 12 /tmp/screen-recording.mp4; fi'
autostart=true
autorestart=false
startretries=0
stderr_logfile=/var/log/ffmpeg.err.log
stdout_logfile=/var/log/ffmpeg.out.log
