# services/openbox.ini
[program:openbox]
command=/usr/bin/openbox
autostart=true
autorestart=true
stderr_logfile=/var/log/openbox.err.log
stdout_logfile=/var/log/openbox.out.log
