# services/vncserver.ini
[program:vncserver]
command=sh -c 'set -x; if [ "$DEBUG" == "true" ]; then /usr/bin/x11vnc; fi'
autostart=true
autorestart=false
startretries=0
stderr_logfile=/var/log/vncserver.err.log
stdout_logfile=/var/log/vncserver.out.log
