# services/vncserver.ini
[program:vncserver]
command=sh -c 'set -x; sleep 5; if [ "$DEBUG" == "true" ]; then /usr/bin/x11vnc -forever; fi'
autostart=true
autorestart=false
startretries=3
stderr_logfile=/var/log/vncserver.err.log
stdout_logfile=/var/log/vncserver.out.log
