# services/xeyes.ini
[program:xeyes]
command=sh -c 'set -x; if [ "$DEBUG" == "true" ]; then /usr/bin/xeyes; fi'
autostart=true
autorestart=false
startretries=3
stderr_logfile=/var/log/xeyes.err.log
stdout_logfile=/var/log/xeyes.out.log
