# services/xvfb.ini
[program:xvfb]
command=/usr/bin/Xvfb %(ENV_DISPLAY)s -screen 0 1920x1080x24+32 -ac +extension RANDR
autostart=true
autorestart=true
stderr_logfile=/var/log/xvfb.err.log
stdout_logfile=/var/log/xvfb.out.log
